# orion-ld.service

[Unit]
Description=FIWARE Context Broker Orion LD

[Service]
ExecStart=/home/ubuntu/fiware.setup/run.sh

[Install]
WantedBy=multi-user.target